# bandersnatch/pypimirrord.conf
[program:tengine]
command=/opt/tengine/sbin/nginx
autostart = true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:bandersnatch]
command=/bin/bash /etc/run.sh
autostart = true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/bandersnatch.log
stdout_logfile_maxbytes=50MB
